# systemd/remote-desktop.service
[Unit]
Description=WayScreenPass Headless Remote Desktop Server
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/remote-desktop-server start
Restart=on-failure
RestartSec=5

# Environment
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
